# mybashd.conf
[program:my_bash_daemon]
command = /usr/local/bin/my_bash_daemon.sh
